# .docker/etc/supervisor.d/api.ini
[program:api]
startsecs = 0
autorestart = true
command=node /var/www/api/src/server.js
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
